# clients/windows-tests/osquery.conf
{
  "options": {
    "config_plugin": "filesystem",
    "logger_plugin": "filesystem",
    "logger_path": "C:\\ProgramData\\osquery\\logs",
    "database_path": "C:\\ProgramData\\osquery\\osquery.db",
    "verbose": false,
    "worker_threads": 2,
    "enable_monitor": true,
    "host_identifier": "hostname",
    "schedule_splay_percent": 10
  },
  "schedule": {
    "system_info": {
      "query": "SELECT hostname, cpu_brand, physical_memory, hardware_vendor, hardware_model FROM system_info;",
      "interval": 3600,
      "description": "Basic system information"
    },
    "os_version": {
      "query": "SELECT name, version, build, platform, platform_like, arch FROM os_version;",
      "interval": 3600,
      "description": "Operating system version information"
    },
    "installed_applications": {
      "query": "SELECT name, version, publisher, install_date, install_location FROM programs;",
      "interval": 3600,
      "description": "Installed Windows applications"
    },
    "running_processes": {
      "query": "SELECT pid, name, path, cmdline, parent, start_time FROM processes WHERE name != '';",
      "interval": 300,
      "description": "Currently running processes"
    },
    "network_interfaces": {
      "query": "SELECT interface, address, mask, broadcast, point_to_point FROM interface_addresses WHERE interface NOT LIKE 'Loopback%';",
      "interval": 1800,
      "description": "Network interface configuration"
    },
    "security_features": {
      "query": "SELECT * FROM windows_security_center;",
      "interval": 1800,
      "description": "Windows security center status"
    },
    "bitlocker_status": {
      "query": "SELECT drive_letter, protection_status, lock_status, encryption_method FROM bitlocker_info;",
      "interval": 3600,
      "description": "BitLocker encryption status"
    },
    "windows_updates": {
      "query": "SELECT title, description, date FROM windows_updates WHERE date > datetime('now', '-7 days');",
      "interval": 3600,
      "description": "Recent Windows updates"
    }
  },
  "packs": {
    "incident-response": {
      "queries": {
        "listening_ports": {
          "query": "SELECT pid, port, protocol, family, address FROM listening_ports;",
          "interval": 600,
          "description": "Applications listening on network ports"
        },
        "startup_items": {
          "query": "SELECT name, path, args, source FROM startup_items;",
          "interval": 3600,
          "description": "Programs that start automatically"
        },
        "logged_in_users": {
          "query": "SELECT user, tty, host, time, pid FROM logged_in_users;",
          "interval": 300,
          "description": "Currently logged in users"
        }
      }
    }
  },
  "decorators": {
    "load": [
      "SELECT uuid AS host_uuid FROM system_info;",
      "SELECT user AS username FROM logged_in_users WHERE user <> '' ORDER BY time DESC LIMIT 1;"
    ]
  }
}
